# etc/systemd/system/onlyoffice-fix.service
[Unit]
Description=Auto-fix OnlyOffice for Nextcloud AIO
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/onlyoffice-auto-fix.sh
